# cron.d
*/5 * * * * /usr/bin/php /var/www/html/update.php